# templates/default/exhibitor.upstart.conf.erb
#!upstart
description "Exhibitor Server"

env USER=<%= @user %>

start on startup
stop on shutdown

respawn

exec sudo -u $USER bash -c "java -Dlog4j.configuration=file://<%= @log4j_props %> -jar <%= @jar %> \
    <%= @opts.keys.sort.map { |k| "--#{k.to_s} #{@opts[k]}" }.join(" ") %> <%= @exec_output %>"

post-start script
  <%= @check_script %> || { stop ; exit 1; }
end script
