# templates/default/init.d-zookeeper.erb
#!/bin/bash
#
# Generated by Chef to <%= node['hostname'] %>
#
# chkconfig: 2345 90 10
# description: zookeeper

SERVICE_NAME="zookeeper"
DIR="<%= node['zookeeper']['install_dir'] %>/zookeeper-<%= node['zookeeper']['version'] %>/bin/"
START="zkServer.sh start"
STOP="zkServer.sh stop"
STATUS="zkServer.sh status"

ZOOCFG="<%= node['zookeeper']['conf_file'] %>"

SERVICE_USER="zookeeper"
SHUTDOWN_WAIT=30
RETVAL="0"

start() {
    if [ `/usr/bin/whoami` != $SERVICE_USER ] && \
       [ `/usr/bin/whoami` != "root" ] ; then
        echo -e "\nNOK! You must be $SERVICE_USER or root to run start!!"
        exit 2
    fi

    echo -e "\nStarting $SERVICE_NAME KEY $KEY ..."
    if [[ `/usr/bin/whoami` == $SERVICE_USER ]]; then
        cd $DIR
        ./$START
    else
        echo "Starting $SERVICE_NAME with $SERVICE_USER..."
        cd $DIR
        su $SERVICE_USER ./$START
    fi
}

stop() { 
    if [ `/usr/bin/whoami` != $SERVICE_USER ] && \
       [ `/usr/bin/whoami` != "root" ] ; then
        echo -e "\nNOK! You must be $SERVICE_USER or root to run start!!"
        exit 2
    fi

    echo -e "\nStarting $SERVICE_NAME KEY $KEY ..."
    if [[ `/usr/bin/whoami` == $SERVICE_USER ]]; then
        cd $DIR
        ./$STOP
    else
        echo "Starting $SERVICE_NAME with $SERVICE_USER..."
        cd $DIR
        su $SERVICE_USER ./$STOP
    fi

}

status() {
  cd $DIR
  ./$STATUS

}

#Body main
case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  status)
    status
    ;;
  restart)
    stop
    sleep 3
    start
    ;;
*)
    echo $"Usage: $0 {start|stop|restart|status}"
    exit 1
esac
exit 0
