# templates/default/zookeeper.sysv.erb
#!/bin/bash
#
# Generated by Chef to <%= node['hostname'] %>
#
# chkconfig: 2345 90 10
# description: zookeeper

[ -r /etc/default/zookeeper ] && . /etc/default/zookeeper

[ -d $ZOO_LOG_DIR ] || mkdir -p $ZOO_LOG_DIR
chown $USER:$GROUP $ZOO_LOG_DIR

DIR="<%= node['zookeeper']['install_dir'] %>/zookeeper-<%= node['zookeeper']['version'] %>/bin/"
START="zkServer.sh start"
STOP="zkServer.sh stop"
STATUS="zkServer.sh status"
CURRENT_USER="$(/usr/bin/whoami)"
SHUTDOWN_WAIT=30
RETVAL="0"

start() {
    if [ $CURRENT_USER != $USER ] && \
       [ $CURRENT_USER != "root" ] ; then
        echo -e "\nNOK! You must be $USER or root to run start!!"
        exit 2
    fi

    echo -e "\nStarting $NAME KEY $KEY ..."
    if [[ $CURRENT_USER == $USER ]]; then
        cd $DIR
        ./$START
    else
        echo "Starting $NAME with $USER..."
        cd $DIR
        su $USER ./$START
    fi
}

stop() {
    if [ $CURRENT_USER != $USER ] && \
       [ $CURRENT_USER != "root" ] ; then
        echo -e "\nNOK! You must be $USER or root to run stop!!"
        exit 2
    fi

    echo -e "\nStopping $NAME KEY $KEY ..."
    if [[ $CURRENT_USER == $USER ]]; then
        cd $DIR
        ./$STOP
    else
        echo "Stopping $NAME with $USER..."
        cd $DIR
        su $USER ./$STOP
    fi

}

status() {
    cd $DIR
    ./$STATUS
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    status)
        status
        ;;
    restart)
        stop
        sleep 3
        start
        ;;
*)
    echo $"Usage: $0 {start|stop|restart|status}"
    exit 1
esac
exit 0
